# install/platform/unix/systemd/wwivd.service.template
[Unit]
Description=wwivd: WWIV UNIX Daemon
[Service]
WorkingDirectory=${WWIV_DIR}
ExecStart=${WWIV_DIR}/wwivd
EnvironmentFile=/etc/wwiv/config
[Install]
WantedBy=multi-user.target

