# bbs/admin/unix/systemd/wwivd.service.template
[Unit]
Description=wwivd: WWIV UNIX Daemon
[Service]
ExecStart=${WWIV_DIR}/wwivd
EnvironmentFile=/etc/default/wwivd
[Install]
WantedBy=multi-user.target

